# crontab
* * * * * . /etc/environment; /usr/local/bin/php /var/www/html/fetch_data.php >> /var/log/light_fetch.log 2>&1
